# tdi-employee-https.service
[Unit]
Description=TDI Employee Application with HTTPS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/app/tdi-main-app-employee
ExecStart=/usr/bin/npm run start-https
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=HTTPS=true
Environment=PORT=3443

[Install]
WantedBy=multi-user.target